# phc.service
# Systemd service file for PHC (SIMS)
# Installation:
# 1. Copy this file to: /etc/systemd/system/phc.service
# 2. Reload systemd: sudo systemctl daemon-reload
# 3. Enable service: sudo systemctl enable phc
# 4. Start service: sudo systemctl start phc
# 5. Check status: sudo systemctl status phc

[Unit]
Description=PHC (SIMS) Gunicorn Application Server
Documentation=https://github.com/munaimtahir/sims
After=network.target
Wants=network.target

[Service]
Type=notify
User=munaim
Group=munaim
WorkingDirectory=/home/munaim/srv/apps/sims
Environment="PATH=/home/munaim/srv/apps/sims/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=sims_project.settings"

# Load environment variables from .env file
EnvironmentFile=/home/munaim/srv/apps/sims/.env

# Gunicorn execution
ExecStart=/home/munaim/srv/apps/sims/venv/bin/gunicorn \
          --config /home/munaim/srv/apps/sims/gunicorn_phc.conf.py \
          sims_project.wsgi:application

# Reload command (graceful restart)
ExecReload=/bin/kill -s HUP $MAINPID

# Process management
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10
PrivateTmp=true

# Security hardening
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phc

[Install]
WantedBy=multi-user.target
